# tests/docker/walli-api.conf
[program:walli-api]
command=/usr/local/bin/gunicorn -b :8080 --chdir /app app:app
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
